# roles/notifications/templates/systemd/vpn-success-watcher.service.j2
[Unit]
Description=WireGuard Telegram Notifier
After=network.target

[Service]
ExecStart={{ vpn_watcher_script }}
Restart=always

[Install]
WantedBy=multi-user.target
