# roles/mihomo/templates/mihomo.service.j2
[Unit]
Description=Clash Meta Client
After=network.target

[Service]
ExecStart=/usr/local/bin/mihomo -f /etc/clash/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target