# server/notifications/templates/systemd/ssh-success-watcher.service.j2
[Unit]
Description=SSH Success Login Watcher
After=network.target

[Service]
ExecStart={{ ssh_watcher_script }}
Restart=always
User=root

[Install]
WantedBy=multi-user.target

